# debian/shard.service
[Unit]
Description=Shard Distributed Secret Vault
After=network.target

[Service]
Type=simple
User=shard
Group=shard
WorkingDirectory=/opt/shard
ExecStart=/opt/shard/venv/bin/python -m shard --config /etc/shard/config.ini
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/var/log/shard/shard.log
StandardError=append:/var/log/shard/shard.err.log

[Install]
WantedBy=multi-user.target
